# scripts/crontab
#!/bin/bash

cd ~/github.com/NikhilPeri/data-land
python -m dataland.scheduler
